# fydetabduo-post-install/fix-display.service
[Unit]
Description=Reinitialize the display after the system has bootedm
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/usr/bin/fix-display

[Install]
WantedBy=sysinit.target