# edge2-post-install/bluetooth-khadas.service
[Unit]
Description=Khadas Bluetooth service

[Service]
Type=simple
ExecStart=/usr/bin/bluetooth-fix

[Install]
WantedBy=multi-user.target
